# nsmd.service
[Unit]
After=network.service

[Service]
Type=forking
ExecStart=/home/ubuntu/nsm/start.sh

[Install]
WantedBy=default.target